# services/purge-old-artifacts.timer
[Unit]
Description=Hourly purge at minute 9

[Timer]
OnCalendar=*:09
Persistent=true

[Install]
WantedBy=timers.target
